# installation-scripts/mysgw.service
[Unit]
Description=MySensors Gateway daemon
Requires=network.target

[Service]
ExecStart=mysgw -q

[Install]
WantedBy=multi-user.target
